# crontab
@daily find /sum/fakkin/path -type f \( -iname '*.txt' -and ! -iname robots.txt \) -atime +30 -delete
@daily find /same/fakkin/path -type f \( -iname '*.pdf' \) -atime +7 -delete
